# install/lolbot.upstart
description "Lolicit discord bot"

start on (local-filesystems and net-device-up IFACE!=lo)
stop on runlevel [016]

setuid lolbot
setgid lolbot

kill signal SIGINT

chdir /opt/lolbot
exec ./lolbot.py
